# rudder-webapp/debian/postinst
#!/bin/sh
# postinst script for rudder-inventory-ldap
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
		echo "Reloading syslog"
		invoke-rc.d rsyslog reload
  /etc/init.d/apache2 stop
  a2dissite default
  a2enmod dav_fs
  a2enmod rewrite
  a2enmod proxy_http
  # Add right to apache user to access /var/rudder/inventories/incoming
  chmod 751 /var/rudder/inventories
  chown root:www-data /var/rudder/inventories/incoming
  chmod 2770 /var/rudder/inventories/incoming
  chown root:www-data /var/rudder/inventories/accepted-nodes-updates
  chmod 2770 /var/rudder/inventories/accepted-nodes-updates
  chmod 755 -R /opt/rudder/share/tools
  a2ensite rudder-default
  htpasswd -bc /opt/rudder/etc/htpasswd-webdav-initial rudder rudder
  htpasswd -bc /opt/rudder/etc/htpasswd-webdav rudder rudder

  echo "(Re-)starting Apache HTTPd"
  /etc/init.d/apache2 restart

  # Run any upgrades
  # Note this must happen *before* creating the technique store, as it was moved in version 2.3.2
  # and creating it manually would break the upgrade logic
  /opt/rudder/bin/rudder-upgrade

	# Create and populate technique store
	if [ ! -d /var/rudder/configuration-repository ]; then mkdir -p /var/rudder/configuration-repository; fi
	if [ ! -d /var/rudder/configuration-repository/shared-files ]; then mkdir -p /var/rudder/configuration-repository/shared-files; fi
	if [ ! -d /var/rudder/configuration-repository/techniques ]; then
		cp -a /opt/rudder/share/techniques /var/rudder/configuration-repository/
	fi

	# Only for Ubuntu:
	## Change rsyslog port number since Ubuntu 12.04 doesn't allow to use standard
	## rsyslog port number (https://bugs.launchpad.net/ubuntu/+source/rsyslog/+bug/789174)
	CHECK_DIST=`/usr/bin/lsb_release -c |cut -d ":" -f2 | sed "s/\s*//"`
	CHECK_RSYSLOG_PORT=`grep -E "^rudder.syslog.port=[0-9]+$" /tmp/rudder-web.properties | cut -d "=" -f2`
	if [ "z${CHECK_DIST}" = "zprecise" -a ${CHECK_RSYSLOG_PORT} -ne "5514" ]; then
		echo "Since Ubuntu 12.04, rsyslog port number used by Rudder  need to be changed"
		sed -i "s/^\(rudder.syslog.port=\)514$/\15514/" /opt/rudder/etc/rudder-web.properties
		echo "rsyslog port number changed to 5514"
	fi


  # Restart the webapp
  echo "Restarting Rudder webapp and inventory-endpoint"
  /etc/init.d/jetty restart

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
