# rudder-webapp/debian/postinst
#!/bin/sh
# postinst script for rudder-inventory-ldap
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# Variables
VAR_RUDDER="/var/rudder"
PT_DIR="${VAR_RUDDER}/configuration-repository/policy-templates"

case "$1" in
    configure)
		echo "Reloading syslog"
		invoke-rc.d rsyslog reload
  /etc/init.d/apache2 stop
  a2dissite default
  a2enmod dav_fs
  a2enmod rewrite
  a2enmod proxy_http
  # Add right to apache user to access /var/rudder/inventories/incoming
  chmod 751 /var/rudder/inventories
  chown root:www-data /var/rudder/inventories/incoming
  chmod 2770 /var/rudder/inventories/incoming
  chown root:www-data /var/rudder/inventories/accepted-nodes-updates
  chmod 2770 /var/rudder/inventories/accepted-nodes-updates
  chmod 755 -R /opt/rudder/share/tools
  a2ensite rudder-default
  htpasswd -bc /opt/rudder/etc/htpasswd-webdav-initial rudder rudder
  htpasswd -bc /opt/rudder/etc/htpasswd-webdav rudder rudder

  echo "(Re-)starting Apache HTTPd"
  /etc/init.d/apache2 restart

  # Run any upgrades
  # Note this must happen *before* creating the technique store, as it was moved in version 2.3.2
  # and creating it manually would break the upgrade logic
  /opt/rudder/bin/rudder-upgrade

	# Create and populate technique store
	if [ ! -d /var/rudder/configuration-repository ]; then mkdir -p /var/rudder/configuration-repository; fi
	if [ ! -d /var/rudder/configuration-repository/shared-files ]; then mkdir -p /var/rudder/configuration-repository/shared-files; fi
	if [ ! -d /var/rudder/configuration-repository/techniques ]; then
		cp -a /opt/rudder/share/techniques /var/rudder/configuration-repository/
	fi


	# Migration of PT 'Set the permissions of files' (Ensure that all actions below won't happen if migration has already made)
	if [ ! -f ${PT_DIR}/fileConfiguration/fileSecurity/filesPermissions/1.0/policy.xml ]; then
		## Commit all modifications before migration
		cd ${PT_DIR} && git add . && git add -u && git commit -am "Committing all pending policy template changes for automatic migration of the policy template from ${PT_DIR}/fileConfiguration/security/ to ${PT_DIR}/fileConfiguration/fileSecurity/" || true
		## Create right folder if it doesn't exist
		if [ ! -d ${PT_DIR}/fileConfiguration/fileSecurity/ ]; then
			mkdir -p "${PT_DIR}/fileConfiguration/fileSecurity"
			echo "${PT_DIR}/fileConfiguration/fileSecurity has been created"
		else
			echo "${PT_DIR}/fileConfiguration/fileSecurity already exists"
		fi

		if [ -d ${PT_DIR}/fileConfiguration/security/ ]; then
		## Check that filePermissions.st located in fileConfiguration/security/ is not duplicated and in the right folder
			if [ -d ${PT_DIR}/fileConfiguration/security/filesPermissions/ ]; then
				echo "The Policy Template 'Set the permissions of files' is not correctly located"
				cd ${PT_DIR} && git mv fileConfiguration/security/* fileConfiguration/fileSecurity/
				cd ${PT_DIR} && git commit -m "Correct Policy Template 'Set the permissions of files' location"
				echo "The location of the Policy Template 'Set the permissions of files' is now correct"
			fi
			## Remove the folder which should contain no more files or folder
			rm -rf ${PT_DIR}/fileConfiguration/security/ # Not using git since it can't manage folder without file
			echo  "${PT_DIR}/fileConfiguration/security/ has been removed"
		fi
	fi

  # Restart the webapp
  echo "Restarting Rudder webapp and inventory-endpoint"
  /etc/init.d/jetty restart

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
