# rudder-webapp/debian/postinst
#!/bin/sh
# postinst script for rudder-inventory-ldap
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
		invoke-rc.d rsyslog reload
  /etc/init.d/apache2 stop
  a2dissite default
  a2enmod dav_fs
  a2enmod rewrite
  a2enmod proxy_http
  chown root:www-data /var/rudder/inventories/incoming
  chmod 2770 /var/rudder/inventories/incoming
  chown root:www-data /var/rudder/inventories/accepted-nodes-updates
  chmod 2770 /var/rudder/inventories/accepted-nodes-updates
  chmod 755 -R /opt/rudder/share/tools
  a2ensite rudder-default
  htpasswd -bc /opt/rudder/etc/htpasswd-webdav rudder rudder

  /etc/init.d/apache2 start

  # Run any upgrades
  # Note this must happen *before* creating the technique store, as it was moved in version 2.3.2
  # and creating it manually would break the upgrade logic
  /opt/rudder/bin/rudder-upgrade

	# Create and populate technique store
	if [ ! -d /var/rudder/configuration-repository ]; then mkdir -p /var/rudder/configuration-repository; fi
	if [ ! -d /var/rudder/configuration-repository/shared-files ]; then mkdir -p /var/rudder/configuration-repository/shared-files; fi
	if [ ! -d /var/rudder/configuration-repository/techniques ]; then
		cp -a /opt/rudder/share/techniques /var/rudder/configuration-repository/
	fi

  # Restart the webapp
  echo "Restarting Rudder webapp and inventory-endpoint"
  /etc/init.d/jetty restart

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
