# rudder-agent/solaris/rudder-smf.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle
  SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<service_bundle type="manifest" name="rudder-agent">
    <service version="1" type="service" name="application/rudder-agent">
        <!--
            The following dependency keeps us from starting until the
            multi-user milestone is reached.
        -->
        <dependency restart_on="none" type="service"
            name="multi_user_dependency" grouping="require_all">
            <service_fmri value="svc:/milestone/multi-user"/>
        </dependency>

        <exec_method timeout_seconds="60" type="method" name="start"
            exec=":true"/>
        <exec_method timeout_seconds="60" type="method" name="stop"
            exec=":true"/>

        <instance enabled="true" name="default"/>

        <template>
            <common_name>
                <loctext xml:lang="C">
                        rudder-agent
                </loctext>
            </common_name>
            <description>
                <loctext xml:lang="C">
                        The rudder-agent service.
                </loctext>
            </description>
        </template>
    </service>

    <service version="1" type="service" name="application/rudder-cf-execd">
        <dependency restart_on="restart" type="service"
            name="agent_dependency" grouping="require_all">
            <service_fmri value="svc:/application/rudder-agent"/>
        </dependency>

        <exec_method timeout_seconds="60" type="method" name="start"
            exec="/opt/rudder/bin/cf-execd %{config/verbosity_option}"/>
        <exec_method timeout_seconds="60" type="method" name="stop"
            exec=":kill"/>
        <exec_method timeout_seconds="60" type="method" name="refresh"
            exec=":kill -HUP"/>

        <instance enabled="true" name="default">
                <property_group type="application" name="config">
                    <propval type="astring" name="verbosity_option" value=""/>
                </property_group>
        </instance>

        <template>
            <common_name>
                <loctext xml:lang="C">
                        rudder-cf-execd
                </loctext>
            </common_name>
            <description>
                <loctext xml:lang="C">
                        CFEngine Execution Scheduler
                </loctext>
            </description>
        </template>
    </service>

    <service version="1" type="service" name="application/rudder-cf-serverd">
        <dependency restart_on="restart" type="service"
            name="agent_dependency2" grouping="require_all">
            <service_fmri value="svc:/application/rudder-agent"/>
        </dependency>

        <exec_method timeout_seconds="60" type="method" name="start"
            exec="/opt/rudder/bin/cf-serverd %{config/verbosity_option}"/>
        <exec_method timeout_seconds="60" type="method" name="stop"
            exec=":kill"/>
        <exec_method timeout_seconds="60" type="method" name="refresh"
            exec=":kill -HUP"/>

        <instance enabled="true" name="default">
                <property_group type="application" name="config">
                    <propval type="astring" name="verbosity_option" value=""/>
                </property_group>
        </instance>

        <template>
            <common_name>
                <loctext xml:lang="C">
                        rudder-cf-serverd
                </loctext>
            </common_name>
            <description>
                <loctext xml:lang="C">
                        CFEngine File Server
                </loctext>
            </description>
        </template>
    </service>
</service_bundle>
